# apps/agent/deploy/systemd/proxmox-backup-orchestrator-pbs-agent-api.service
[Unit]
Description=Proxmox Backup Orchestrator PBS execution agent API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/proxmox-backup-orchestrator-agent
EnvironmentFile=/opt/proxmox-backup-orchestrator-agent/.env
ExecStart=/opt/proxmox-backup-orchestrator-agent/.venv/bin/python -m agent.main serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
